# ops/test-bbsd.service
[Unit]
Description=test-bbs daemon (bbsd)
After=network.target

[Service]
Type=simple
User=bbsd
Group=bbs
WorkingDirectory=/opt/test-bbs
Environment=REPO_DIR=/opt/test-bbs
Environment=BBS_SOCKET_PATH=/run/test-bbs/bbsd.sock
Environment=BBS_DB_PATH=/var/lib/test-bbs/bbsd.sqlite3
ExecStart=/opt/test-bbs/ops/bbsd-run.sh
Restart=on-failure
RuntimeDirectory=test-bbs
RuntimeDirectoryMode=0770
UMask=0007

[Install]
WantedBy=multi-user.target

